# systemd/webrtc-sender-rov.service
[Unit]
Description=WebRTC sender (main camera)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=REPO_DIR=/opt/mantaray/Mantaray-on-IaC
WorkingDirectory=${REPO_DIR}
ExecStart=/usr/bin/python3 ${REPO_DIR}/src/webrtc/webrtc_sender.py --device /dev/video0 --signaling ws://127.0.0.1:9000?stream=rov --width 1280 --height 720 --fps 30 --bitrate 2500 --key-int 30
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
